# gunicorn.service
[Unit]
Description=Gunicorn daemon for ML Interview Notes
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/ml-interview-notes
Environment="PATH=/home/ubuntu/ml-interview-notes/venv/bin"
ExecStart=/home/ubuntu/ml-interview-notes/venv/bin/gunicorn --config gunicorn.conf.py ml_interview_notes.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target 